# eredes_meter/rootfs/etc/services.d/eredes_meter/run
#!/usr/bin/with-contenv bashio
# Start the E-REDES meter scraper service
exec /usr/bin/eredes_meter.py
